# com.geomet.copper-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.geomet.copper-dashboard</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/richardgoldberg/geomet-copper-dashboard/app.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/richardgoldberg/geomet-copper-dashboard</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/richardgoldberg/geomet-copper-dashboard/dashboard.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/richardgoldberg/geomet-copper-dashboard/dashboard_err.log</string>
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
